# phoenix_v5.service
[Unit]
Description=Phoenix v5.1.0c Trading Bot
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/phoenix_v5_run
ExecStartPre=/bin/sleep 5
ExecStart=/root/phoenix_v5_run/venv/bin/python main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# 환경변수
Environment=PYTHONUNBUFFERED=1

# 메모리 제한 (선택)
# MemoryLimit=1G

[Install]
WantedBy=multi-user.target
